# base/koolna-base/etc/profile.d/koolna-paths.sh
export PATH="$HOME/.local/bin:$HOME/.local/share/mise/shims:/usr/local/bin:$PATH"
